# app/Supervisor/reverb.ini
[program:PROJECT-reverb]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/virtual/USER/PROJECT/current/artisan reverb:start
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=USER
numprocs=1
redirect_stderr=true
stdout_logfile=/home/USER/reverb.log
stopwaitsecs=3600
startsecs=30
